# users/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:rabbitmq-consumer]
command=php bin/console doctrine:migrations:migrate --no-interaction
